# agent_controller_linux/nsight-agent.service
[Unit]
Description=NSight2.0 Agent
After=network.target

[Service]
Type=simple
ExecStart=PYTHON AGENT

[Install]
WantedBy=multi-user.target

